# scripts/beesd@.service.in
[Unit]
Description=Bees (%i)
Documentation=https://github.com/Zygo/bees
After=local-fs.target systemd-udev-settle.service
Before=umount.target

[Service]
Type=simple
ExecStart=/usr/sbin/beesd --no-timestamps %i

CPUSchedulingPolicy=batch
CPUWeight=12
IOSchedulingClass=idle
IOSchedulingPriority=7
IOWeight=10
KillMode=control-group
KillSignal=SIGTERM
MemoryAccounting=true
Nice=19
Restart=on-abnormal
RuntimeDirectoryMode=0700
RuntimeDirectory=bees
StartupCPUWeight=25
StartupIOWeight=25

# Give beesd time for ZRAM teardown + flush
TimeoutStopSec=90

# Hide other users' process in /proc/
ProtectProc=invisible

# Mount / as read-only, but allow ZRAM sysfs
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateNetwork=true
PrivateIPC=true
ProtectHostname=true
ProtectKernelTunables=false
ProtectControlGroups=false

# Allow ZRAM control + devices
ReadWritePaths=/sys /dev

# Capabilities
AmbientCapabilities=CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER CAP_SYS_ADMIN
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
